# systemd/codex-usage.timer
[Unit]
Description=Run ChatGPT Codex usage analytics poll every five minutes
Documentation=man:systemd.timer(5)

[Timer]
OnBootSec=1min
OnUnitActiveSec=5min
AccuracySec=30s
Persistent=true
Unit=codex-usage.service

[Install]
WantedBy=timers.target
