# systemd/codex-usage.service
[Unit]
Description=Watch ChatGPT Codex usage analytics and block exhausted accounts
Documentation=man:systemd.service(5)

[Service]
Type=simple
ExecStart=%h/.local/bin/codex-usage watchdog --format table
Environment=PYTHONUNBUFFERED=1
TimeoutStartSec=120
RuntimeMaxSec=180
TimeoutStopSec=15
KillMode=mixed
MemoryMax=1G
TasksMax=256
OOMPolicy=kill
Restart=no
UMask=0077
NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectClock=true
ProtectHostname=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=%h/.config/codex-usage %h/.local/share/codex-usage %h/.cache/ms-playwright
RestrictSUIDSGID=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
LockPersonality=true
SystemCallArchitectures=native

[Install]
WantedBy=default.target
